# assignment3/EC2_userdata/pumpkinservice.conf
start on startup
stop on shutdown
setuid ubuntu
respawn
respawn limit 20 5

script
exec /home/ubuntu/run_pumpkin.sh
end script